# systemd/instrumentation.target
[Unit]
Description=Nottinghack Instrumentation
#After=Multi-User

[Install]
WantedBy=multi-user.target


